# deploy/systemd/ferry-services-timetable-document-fetcher.service
[Unit]
Description=Scottish Ferry Services timetable document fetcher
After=network-online.target
Wants=network-online.target
PartOf=ferry-services.target

[Service]
User=stefanchurch
Type=oneshot
WorkingDirectory=/home/stefanchurch/ferry-services-server-v3
EnvironmentFile=-/home/stefanchurch/ferry-services-server-v3/.env
Environment=NODE_ENV=production
ExecStart=/usr/bin/npm run fetch:timetable-documents
